# deb_repo/thirdreality-hacore/matter-server.service
[Unit]
Description=Matter Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=PATH=/srv/homeassistant/bin:/usr/local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStartPre=/bin/sleep 5
ExecStart=/srv/homeassistant/bin/matter-server --bluetooth-adapter 0 --storage-path /var/lib/homeassistant/matter_server --paa-root-cert-dir /var/lib/homeassistant/matter_server/credentials --log-file /var/lib/homeassistant/matter_server/matter_server.log

[Install]
WantedBy=multi-user.target
